# packages/bbperf/bbperf-0.0.10.tar.gz/bbperf-0.0.10/bbperf.service
[Unit]
Description=bbperf
After=syslog.target network.target

[Service]
Type=simple
User=mfreemon
ExecStart=/home/mfreemon/.venv-bbperf/bin/bbperf -s
Restart=always

[Install]
WantedBy=multi-user.target

